# kindle/etc/crontab/root
# */15 * * * * /usr/sbin/checkpmond
# */15 * * * * /usr/sbin/tinyrot
# */60 * * * * /usr/sbin/loginfo tmpfs
# */60 * * * * /usr/sbin/loginfo localVars
# */60 * * * * /usr/sbin/loginfo memusedump
# */15 * * * * /usr/sbin/loginfo powerdcheck

*/5 6-22 * * * /mnt/us/update.sh
0 23,0,5 * * * /mnt/us/update.sh
